# conf/nginx-supervisor.ini
[supervisord]
nodaemon=true

[program:nginx]
command = /usr/sbin/nginx
user = root
autostart = true
stdout_logfile=/var/log/nginx/access.log
stderr_logfile=/var/log/nginx/error.log

[program:php-fpm]
command = /usr/bin/php-fpm -F
user = root
autostart = true
stderr_logfile=/var/log/php-fpm/error.log
